# recipes-multimedia/pistreamudp/files/systemd-units/pistreamudp.service
[Unit]
Description=PiStreamUDP
After=multi-user-.target

[Service]
Type=simple
ExecStart=/bin/pistreamudp.sh

[Install]
WantedBy=multi-user.target
