# recipes-multimedia/pistream/files/systemd-units/pistream.service
[Unit]
Description=PiStream
After=multi-user-.target

[Service]
Type=simple
ExecStart=/bin/pistream.sh

[Install]
WantedBy=multi-user.target
