# recipes-misc/pi-garage-door/files/systemd-units/pi-garage-door.service
[Unit]
Description=Pi Garage Door Output Pin Initializer

[Service]
Type=oneshot
ExecStart=/bin/pi-garage-door-init.sh
 
[Install]
WantedBy=multi-user.target